# ansible/roles/reptar.prakashsurya.com/files/etc/systemd/system/x11vnc.service
[Unit]
Description=Start x11vnc at startup.
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/x11vnc -auth guess -forever -loop -rfbport 5900

[Install]
WantedBy=multi-user.target
